# gps_conversion_tool/src/gps_conversion_tool.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/run_gps_conversion_tool.sh

[Install]
WantedBy=default.target